# crontab
0 5 * * * root python /app/scripts/update_db.py >> var/log/cron.log 2>1